# azure_foundry_ai/rootfs/etc/services.d/azure_foundry/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Azure Foundry AI service..."
exec python3 /usr/bin/azure_foundry_service.py